# docs/feedgen-fastcgi.service
[Unit]
Description=Feedgen FastCGI
After=syslog.target

[Service]
ExecStart=/bin/bash -l -c "uwsgi --ini /srv/feedgen.hasname.com/uwsgi.ini"
Group=www-data
KillSignal=SIGINT
Restart=always
RuntimeDirectory=uwsgi
Type=simple
User=www-data

[Install]
WantedBy=multi-user.target
